# {{cookiecutter.repo_name}}/menu_qr.conf
[program:menu_qr]
command=docker-compose up
directory=/root/apps/menu-fe-compose
redirect_stderr=true
autostart=true
autorestart=true
priority=10
