# {{cookiecutter.repo_name}}/supervisor.ini
[program:{{repo_name}}]
command=docker-compose -f production.yml up
directory=/home/ec2-user/{{repo_name}}
redirect_stderr=true
autostart=true
autorestart=true
priority=10
