# {{cookiecutter.repo_name}}/supervisor.ini
[program:{{cookiecutter.repo_name}}]
command=sudo docker-compose up
directory=/home/ec2-user/{{cookiecutter.repo_name}}
redirect_stderr=true
autostart=true
autorestart=true
priority=10
